# roles/raspberry/files/keep_connection_alive.service
[Unit]
Description=Keep connection alive

[Service]
Type=simple
ExecStart=/opt/keep_connection_alive.sh
